# usr/share/manjaro-tools/iso-profiles/shared/manjaro/live-overlay/etc_bloqueado/udev/rules.d/10-usb-ro.rules
ACTION=="add|change",SUBSYSTEM=="usb",ATTRS{removable}=="1",RUN+="/sbin/blockdev --setrw /dev/%k"
ACTION=="add|change",SUBSYSTEM=="block",ATTRS{removable}=="0",RUN+="/sbin/blockdev --setro /dev/%k"
